# .github/actions/ansible/etc/profile.d/ansible.sh
ANSIBLE_HOME=/appli/softs/ansible
PYTHONPATH=/appli/softs/ansible/lib:/usr/lib/python3.8/site-packages
PATH=$ANSIBLE_HOME/bin:$PATH

export ANSIBLE_HOME
export PYTHONPATH
export PATH